# com.vynzyx.vencord-reinject.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vynzyx.vencord-reinject</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{SCRIPT}}</string>
    </array>

    <!-- Fire the instant Discord's bundle changes, i.e. an app update swaps app.asar. -->
    <key>WatchPaths</key>
    <array>
        <string>/Applications/Discord.app/Contents/Resources</string>
    </array>

    <!-- Backstops in case a watch event is missed. -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Coalesce rapid bundle changes into at most one run per 10s. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardErrorPath</key>
    <string>{{ERRLOG}}</string>
</dict>
</plist>
